# ci/nodepool/nodepoold.service
[Unit]
Description=nodepoold

[Service]
Type=simple
ExecStart=/bin/nodepoold -d
User=nodepool